# service/linux/slack-status.timer
[Unit]
Description=Hourly update Slack status with current location

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
